# .github/wiki-notebook.service
[Unit]
Description=Wiki Notebook - Self-Hosted Journal and Wiki
Documentation=https://github.com/coreconduit/wiki-notebook
After=network.target

[Service]
Type=notify
User=wiki-notebook
Group=wiki-notebook
WorkingDirectory=/opt/wiki-notebook

# Environment
Environment="FLASK_ENV=production"
Environment="PYTHONUNBUFFERED=1"

# Execution
ExecStart=/opt/wiki-notebook/.venv/bin/python -m wiki_notebook
Restart=on-failure
RestartSec=10s

# Security Hardening
ProtectSystem=strict
ProtectHome=true
NoNewPrivileges=true
PrivateTmp=true
ProtectClock=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
ProtectProc=invisible
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictNamespaces=true

# ReadWrite paths: database directory
ReadWritePaths=/opt/wiki-notebook/data

# Standard output to journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wiki-notebook

[Install]
WantedBy=multi-user.target
