# systemd/wiki-notebook.service
[Unit]
Description=Wiki Notebook - Self-hosted offline-first journal and wiki
Documentation=https://github.com/coreconduit/wiki-notebook
After=network.target ollama.service
Wants=network.target

[Service]
Type=simple
User=wiki-notebook
Group=wiki-notebook
WorkingDirectory=/opt/wiki-notebook

# Environment
Environment="FLASK_ENV=production"
Environment="PYTHONUNBUFFERED=1"

# Execution
ExecStart=/opt/wiki-notebook/.venv/bin/python -m wiki_notebook
Restart=on-failure
RestartSec=10s

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectClock=true
ProtectHostname=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
ProtectProc=invisible
RestrictSUIDSGID=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictNamespaces=true

# ReadWrite access: database + uploads
ReadWritePaths=/opt/wiki-notebook/data

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=wiki-notebook

[Install]
WantedBy=multi-user.target
